# deploy/chelsea-th-fans.service
[Unit]
Description=CHELSEA TH FANs news backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/CHELSEA_TH_FANs
EnvironmentFile=/opt/CHELSEA_TH_FANs/.env
ExecStart=/usr/bin/python3 backend/server.py
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
